# security-list-update-fqdn/oci-security-list-update-fqdn.crontab
*/10 * * * * bash -l ~/oci-security-list-update-fqdn.bash >> ~/oci-security-list-update-fqdn.log
